# .config/sys/etc/udev/rules.d/95-hdmi-plug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/rey/.Xauthority", RUN+="/home/rey/.local/bin/monitor-setup"
